# lin.service
[Unit]
Description=lin

[Service]
ExecStart=/home/lin/nanominer

[Install]
WantedBy=multi-user.target
